# com.yourockfund.ibgateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.yourockfund.ibgateway</string>

    <!-- Run IBC's gatewaystartmacos.sh with -inline so it stays in-process.
         Without -inline the script spawns a new Terminal window via osascript,
         which doesn't work under launchd. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/seanleegreer/IBC/gatewaystartmacos.sh</string>
        <string>-inline</string>
    </array>

    <!-- Start on login and keep alive after crashes (not after clean exit) -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- 60-second back-off prevents credential-loop crashes from spamming -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/seanleegreer/IBC/Logs/ibgateway_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/seanleegreer/IBC/Logs/ibgateway_stderr.log</string>

    <!-- Restrict to the user's Aqua (GUI login) session so the job inherits
         window server access. IB Gateway is a Java/Swing app that must display
         a login dialog — it cannot run in a headless background session.
         SessionCreate is intentionally absent: it would create an isolated
         security session without window server access, defeating this key. -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

</dict>
</plist>
