# com.yourockfund.docker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yourockfund.docker</string>

    <!-- Runs setup_docker.sh --paper on every login: pulls secrets from Keychain,
         writes ephemeral docker/secrets/ files, starts all 4 containers, then
         deletes the secret files. The 30-second sleep gives Rancher Desktop time
         to fully start before the script runs. -->
    <!-- setup_docker.sh installs this to ~/Library/LaunchAgents/ with __PROJ__
         and __HOME__ substituted. -->

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>sleep 60 &amp;&amp; chmod +x __PROJ__/setup_docker.sh &amp;&amp; exec __PROJ__/setup_docker.sh --paper</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJ__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Cover both Homebrew paths and Rancher Desktop's /usr/local/bin symlink. -->
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>300</integer>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/yrvi-autostart.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/yrvi-autostart.log</string>
</dict>
</plist>
